# deploy/skyclaw.service
[Unit]
Description=SkyClaw AI Agent Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=skyclaw
Group=skyclaw
WorkingDirectory=/var/lib/skyclaw
EnvironmentFile=/etc/skyclaw/env
ExecStart=/usr/local/bin/skyclaw start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
